# mender/image-tools/artifact-builder/raspbian-setup/files/first-boot.service
[Unit]
Description=First Boot
DefaultDependencies=no
After=local-fs.target

[Service]
ExecStart=/root/first-boot.sh

[Install]
WantedBy=sysinit.target
