# mender/image-tools/artifact-builder/raspbian-setup/files/network-watchdog.service
[Unit]
Description=Network Watchdog
After=network-online.target

[Service]
ExecStart=/root/network-watchdog.sh

[Install]
WantedBy=multi-user.target
