# mender/image-tools/artifact-builder/raspbian-setup/files/video-updater.service
[Unit]
Description=Video Updater
After=network-online.target

[Service]
ExecStart=/root/video-updater.sh

[Install]
WantedBy=multi-user.target
