# mender/image-tools/artifact-builder/raspbian-setup/files/web-service.conf
[program:web-service-flask-app]
directory=/var/flaskapp/web-service
command=/bin/bash -E -c ./start.sh
autostart=true
autorestart=true
stopsignal=INT
stopasgroup=true
killasgroup=true
user=www-data
