# mender/image-tools/artifact-builder/raspbian-setup/files/network/rfkill-unblock.service
[Unit]
Description=RFKill-Unblock Wifi

[Service]
Type=oneshot
ExecStart=/usr/sbin/rfkill unblock wifi
ExecStop=
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target